# deploy/qubes/db-qube/ob1-db-firewall.service
[Unit]
# Re-applies the tailscale0:5432 nft accept AFTER tailscaled brings the interface
# up. qubes-firewall.service runs before tailscaled, when tailscale0 does not yet
# exist, so the rule from qubes-firewall-user-script would otherwise be lost.
#
# Stash at: /rw/config/ob1-db-firewall.service
# rc.local copies it to /etc/systemd/system/ and enables it each boot (/etc is
# wiped on every AppVM reboot, so it cannot be enabled persistently in place).
Description=OB1 DB qube post-tailscale firewall (tcp/5432 on tailscale0)
After=tailscaled.service
Wants=tailscaled.service

[Service]
Type=oneshot
ExecStart=/rw/config/qubes-firewall-user-script
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
